# airootfs/root/endeavouros-skel-liveuser/src/etc/skel/.config/autostart/theme.desktop
[Desktop Entry]
Exec=/home/liveuser/.theming.sh
Icon=dialog-scripts
Name=Theme
Type=Application
X-KDE-AutostartScript=true
